# systemd/crawler-vagas-dev-13h.service
[Unit]
Description=Crawler Vagas Dev - slot das 13h
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env
ExecStart=/usr/bin/npm run run:cron:13
StandardOutput=append:__PROJECT_DIR__/data/cron.log
StandardError=append:__PROJECT_DIR__/data/cron.log

[Install]
WantedBy=default.target
