# systemd/crawler-vagas-dev-08h.service
[Unit]
Description=Crawler Vagas Dev - slot das 08h
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env
Environment=PATH=__NODE_BIN__:/usr/local/bin:/usr/bin:/bin
ExecStart=__NODE_BIN__/npm run run:cron:8
StandardOutput=append:__PROJECT_DIR__/data/cron.log
StandardError=append:__PROJECT_DIR__/data/cron.log

[Install]
WantedBy=default.target
